# packaging/debian/prerm
#!/bin/sh
set -e

case "$1" in
    remove|deconfigure)
        if command -v systemctl >/dev/null 2>&1; then
            systemctl stop tv-converter.service || true
            systemctl disable tv-converter.service || true
        fi
        ;;
esac

exit 0
