# systemd/tv-converter.service
[Unit]
Description=TV Converter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

WorkingDirectory=/home/plex/tv-converter

ExecStart=/home/plex/tv-converter/venv/bin/python \
    /home/plex/tv-converter/main.py \
    --config /home/plex/tv-converter/config.yaml

ExecReload=/bin/kill -HUP $MAINPID

Restart=always
RestartSec=10

KillSignal=SIGTERM
TimeoutStopSec=30

Environment=PYTHONUNBUFFERED=1

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
