# systemd/tv-converter.service.in
[Unit]
Description=TV Converter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@SERVICE_USER@
Group=@SERVICE_GROUP@
WorkingDirectory=/var/lib/tv-converter
ExecStart=/var/lib/tv-converter/venv/bin/python \
    /var/lib/tv-converter/main.py \
    --config /etc/tv-converter/config.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
KillSignal=SIGTERM
TimeoutStopSec=30
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
